# debian/i2pd.upstart
description "dotnet client daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

# these can be overridden in /etc/init/dotnet.override
env LOGFILE="/var/log/dotnet/dotnet.log"

expect fork

exec /usr/sbin/dotnet --daemon --service --log=file --logfile=$LOGFILE
